# {{cookiecutter.repo_name}}/docker/python/worker.ini
[supervisord]
nodaemon = true
nocleanup = true
logfile = /dev/stdout
logfile_maxbytes = 0

[program:worker]
autorestart = true
command = django-admin rqworker
stopasgroup = true
numprocs = 3
numprocs_start = 1
process_name = %(program_name)s-%(process_num)s
stopsignal = TERM
user = app
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
